# tmux/scheduler/tmux-snapshot.timer
[Unit]
Description=Run tmux snapshot every 5 minutes

[Timer]
OnBootSec=5min
OnUnitActiveSec=5min
AccuracySec=30s
Persistent=false

[Install]
WantedBy=timers.target
